# udev_rules/100-external_display_switch.rules
ACTION=="change", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/minho/.Xauthority", RUN+="/home/minho/.local/bin/display-ctl -a 5"
